# contrib/systemd/systemd.service.template
[Unit]
Description=IntelMQ Bot %i

[Service]
ExecStart={{MODULE}} %i
Restart=on-failure
StandardOutput=journal
